# systemd/go-api.service
[Unit]
Description=Go Backend API Service
After=network.target librechat-backend.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/librechat_me/saas-api/cmd/api
Environment="PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/local/go/bin/go run main.go
Restart=always
RestartSec=10
StandardOutput=append:/home/ec2-user/librechat_me/logs/go-api.log
StandardError=append:/home/ec2-user/librechat_me/logs/go-api-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

